# systemd/modem-routing.service
[Unit]
Description=Modem Multipath Routing Setup
After=network-online.target ModemManager.service NetworkManager.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/root/repo/proxyfarm/scripts/setup_modem_routing.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

# Restart routing if it fails
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
