# systemd/proxyfarm-proxy.service
[Unit]
Description=ProxyFarm HTTP/HTTPS Proxy Server
After=network.target ModemManager.service NetworkManager.service
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/repo/proxyfarm
# Use venv if available, fallback to system python
ExecStart=/bin/bash -c 'source /root/proxyfarm-venv/bin/activate 2>/dev/null && python proxy_full.py || /usr/bin/python3 proxy_full.py'
Restart=always
RestartSec=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=proxyfarm-proxy

# Security (optional)
# NoNewPrivileges=true
# PrivateTmp=true

[Install]
WantedBy=multi-user.target
